# crontab
# Delete old SSL user certs
7 * * * *  find /home/eric/certs/generated -mtime +1 -exec rm {} \;
